# ngrok-ssh.service
[Unit]
Description=SSH Tunnel watcher
After=network.target

[Service]
Type=simple
ExecStart=%h/%N/ngrok-start.sh

[Install]
WantedBy=default.target
